# etc/zigzag-languages.service
[Unit]
Description=Install language packages

[Service]
Type=oneshot
ExecStart=/usr/bin/zigzag-languages -y
